# asterisk_mon_start.conf
# /etc/init script
#
# $Id: asterisk_mon_start.conf 3566 2010-06-29 19:00:01Z marlon $

start on runlevel [2345]
stop on runlevel [!2345]

respawn
exec /usr/local/sbin/asterisk_mon_start.py
